# systemd/alienware-vllm-strong.service
[Unit]
Description=Homelab vLLM strong-model backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/homelab-control/vllm-strong.env
ExecStartPre=/usr/bin/podman rm -f homelab-vllm-strong
ExecStart=/usr/bin/bash -lc '/usr/bin/podman run --name homelab-vllm-strong --replace --device nvidia.com/gpu=all --security-opt=label=disable --ipc=host -p 0.0.0.0:8001:8001 -v %h/.cache/huggingface:/root/.cache/huggingface:Z --env-file %h/.config/homelab-control/vllm-strong.env docker.io/vllm/vllm-openai:v0.19.1 --host 0.0.0.0 --port 8001 --model "$VLLM_STRONG_MODEL" --served-model-name "$VLLM_STRONG_SERVED_MODEL_NAME" --dtype "$VLLM_STRONG_DTYPE" --max-model-len "$VLLM_STRONG_MAX_MODEL_LEN" --gpu-memory-utilization "$VLLM_STRONG_GPU_MEMORY_UTILIZATION" --kv-cache-dtype "$VLLM_STRONG_KV_CACHE_DTYPE" --max-num-batched-tokens "$VLLM_STRONG_MAX_NUM_BATCHED_TOKENS" --max-num-seqs "$VLLM_STRONG_MAX_NUM_SEQS" --enable-prefix-caching --enable-chunked-prefill --api-key "$VLLM_STRONG_API_KEY" $VLLM_STRONG_EXTRA_ARGS'
ExecStop=/usr/bin/podman stop -t 20 homelab-vllm-strong
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
