# systemd/alienware-agent-platform-report.timer
[Unit]
Description=Refresh homelab agent platform status snapshot every 5 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
Unit=alienware-agent-platform-report.service

[Install]
WantedBy=timers.target
