# systemd/alienware-health-monitor.service
[Unit]
Description=Homelab health monitor pass (audit chains, timers, services, endpoints, restic)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/homelab-control/health-monitor.env
WorkingDirectory=%h/git/homelab-control
ExecStart=%h/git/homelab-control/.venv/bin/python -m apps.health_monitor
Nice=10
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3min

[Install]
WantedBy=default.target
