# systemd/alienware-executive-chat.service
[Unit]
Description=Homelab executive-assistant local chat UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/homelab-control/agent-executive.env
EnvironmentFile=-%h/.config/homelab-control/agent-executive-chat.env
ExecStart=/usr/bin/bash -lc 'cd "${HOMELAB_CONTROL_ROOT:-%h/git/homelab-control}" && exec python3 apps/executive_agent/chat_server.py --state-dir %h/.local/state/homelab-control/agent-executive'
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
