# systemd/alienware-executive-weekly-review.service
[Unit]
Description=Homelab executive-assistant weekly operating review
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/homelab-control/agent-executive.env
ExecStart=/usr/bin/bash -lc 'cd "${HOMELAB_CONTROL_ROOT:-%h/git/homelab-control}" && exec python3 apps/executive_agent/main.py weekly-review --state-dir %h/.local/state/homelab-control/agent-executive --output %h/.local/state/homelab-control/agent-executive/weekly-review.json'
