# systemd/alienware-agent-homelab-discord.service
[Unit]
Description=Homelab agent homelab Discord bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=%h/.config/homelab-control/agent-homelab.env
EnvironmentFile=%h/.config/homelab-control/agent-homelab-discord.env
ExecStart=/usr/bin/bash -lc 'cd "%h/git/homelab-control" && exec python3 apps/author_agent/discord_bot.py --state-dir %h/.local/state/homelab-control/agent-homelab'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
