# systemd/alienware-author-agent.service
[Unit]
Description=Homelab author-agent queue worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/homelab-control/agent-homelab.env
ExecStartPre=/usr/bin/bash -lc 'mkdir -p %h/.local/state/homelab-control/agent-homelab/{inbox,processing,done,failed}'
ExecStart=/usr/bin/bash -lc 'cd "${HOMELAB_CONTROL_ROOT:-%h/homelab-control}" && exec python3 apps/author_agent/main.py worker --queue-dir %h/.local/state/homelab-control/agent-homelab --heartbeat %h/.local/state/homelab-control/agent-homelab/heartbeat.json'
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
