# systemd/alienware-review-agent.service
[Unit]
Description=Homelab review-agent queue worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/homelab-control/agent-review.env
WorkingDirectory=%h/homelab-control
ExecStartPre=/usr/bin/bash -lc 'mkdir -p %h/.local/state/homelab-control/agent-review/{inbox,processing,done,failed}'
ExecStart=/usr/bin/env python3 %h/homelab-control/apps/review_agent/main.py worker --queue-dir %h/.local/state/homelab-control/agent-review --heartbeat %h/.local/state/homelab-control/agent-review/heartbeat.json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
