# systemd/homelab-operator.service
[Unit]
Description=Homelab operator reports

[Service]
Type=oneshot
WorkingDirectory=/opt/homelab-control
ExecStart=/usr/bin/env python3 /opt/homelab-control/apps/homelab_operator/main.py observability-report
