# systemd/homelab-operator.service
[Unit]
Description=Homelab inventory memory sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/homelab-operator
ExecStartPre=/usr/bin/env bash -lc 'cd "${HOMELAB_CONTROL_ROOT:-/opt/homelab-control}" && ./scripts/render-operator-env.sh'
ExecStart=/usr/bin/env bash -lc 'cd "${HOMELAB_CONTROL_ROOT:-/opt/homelab-control}" && exec python3 apps/homelab_operator/main.py inventory-memory-sync'
